# tools/scripts/crontab
*/5 * * * * /data/homelab/tools/duckdns/duck.sh >/dev/null 2>&1